# containers/atvloadly/atvloadly.container
# atvloadly — Apple TV IPA sideloading.
# IMPORTANT: SeccompProfile=unconfined and the host dbus/avahi mounts are
# required for device pairing — do not remove or tighten them.
# Install: cp atvloadly.container ~/.config/containers/systemd/
#          systemctl --user daemon-reload && systemctl --user start atvloadly
[Unit]
Description=atvloadly — Apple TV IPA sideloading
Wants=network-online.target
After=network-online.target

[Container]
Image=docker.io/bitxeno/atvloadly:latest
ContainerName=atvloadly
HostName=atvloadly
AutoUpdate=registry

# Scheduler/cron runs on local time, not UTC.
Environment=TZ=Europe/Amsterdam

# --- security DEVIATION (intentional): unconfined seccomp for usbmuxd/device
#     access; caps intentionally NOT dropped and NoNewPrivileges NOT set, like
#     the Docker stack — tightening them breaks pairing. ---
SeccompProfile=unconfined

# --- UID mapping: /etc/atvloadly is owned by your real host user, mode 700
#     (owner-only — verified live). Docker accessed it as real root (bypasses
#     all permission checks); rootless Podman never gets real root, so the
#     container's OWN identity must literally be your UID. keep-id + User=
#     do that — same pattern as jellyfin. PUID/PGID don't exist for this
#     image (it has no such switching logic), so without this it would run
#     as a root that doesn't match the directory's real owner and get
#     "Permission denied" on every read/write, same class of bug as
#     qBittorrent's original issue. host dbus/avahi sockets are world
#     read-write (verified, mode 666) — no related fix needed for those.
UserNS=keep-id
User=1000:1000
# The app binds directly to port 80 INSIDE the container (the Docker stack
# never noticed because Docker ran it as real root). As non-root UID 1000,
# binding a port <1024 needs this capability — verified live: without it,
# "failed to listen: listen tcp4 0.0.0.0:80: bind: permission denied" and an
# immediate crash loop, regardless of which HOST port is published.
AddCapability=NET_BIND_SERVICE

# --- networking: host, not bridge (verified — the ORIGINAL Docker stack used
#     network_mode: host, which our first conversion missed; it defaulted to
#     Podman's isolated bridge network instead). This matters for more than
#     the web UI: Apple's device-pairing protocol needs broad LAN reachability
#     that a NAT'd bridge blocks — avahi discovery still worked (it queries
#     the host's avahi socket directly), but the app showed an empty device
#     list and no live connection until switched to host networking.
#     PublishPort= is meaningless under Network=host (no port mapping occurs;
#     the app binds directly to the HOST's port 80 — the web UI moves from
#     :5533 to :80 as a result). Rootless Podman supports host networking
#     fine — it reuses the host's existing netns, no extra privilege needed.
Network=host

# --- storage ---
# /data = ALL persistent state: pairing, Apple ID session, app.db, settings (backup target)
Volume=/etc/atvloadly:/data
Volume=/var/run/dbus:/var/run/dbus
Volume=/var/run/avahi-daemon:/var/run/avahi-daemon
Volume=/etc/localtime:/etc/localtime:ro

[Service]
Restart=always

[Install]
WantedBy=default.target
