# containers/qbittorrent/qbittorrent.container
# qBittorrent — rides entirely on gluetun's network namespace. No network stack
# of its own: if gluetun stops or the tunnel drops, qBittorrent loses network
# (the kill switch). Because it shares gluetun's netns, it gets NO PublishPort of
# its own — its WebUI is published on gluetun.container instead.
[Unit]
Description=qBittorrent (in gluetun's netns)
# Start after gluetun and require it — sharing its netns depends on it existing.
Requires=gluetun.service
After=gluetun.service

[Container]
Image=lscr.io/linuxserver/qbittorrent:latest
ContainerName=qBittorrent
AutoUpdate=registry

# Share gluetun's netns (= Docker's `network_mode: service:gluetun`).
Network=container:gluetun

# --- security: linuxserver/s6 image; caps intentionally NOT dropped (unverified
#     against this image's init, like the Docker stack). NoNewPrivileges kept. ---
NoNewPrivileges=true
PodmanArgs=--memory=1g --pids-limit=400

# --- UID mapping: NOT keep-id here. The linuxserver/s6 image remaps to
#     PUID/PGID internally on its own (that's what those .env vars are for) —
#     and rootless Podman does not allow UserNS=keep-id together with
#     Network=container:<other> (joining a netns can't also set up a separate
#     user namespace). keep-id is only needed for images with no PUID/PGID
#     logic of their own (e.g. jellyfin).
EnvironmentFile=%h/containers/qbittorrent/.env
Environment=WEBUI_PORT=8080

# --- storage: app config + your downloads disk (DOWNLOADS_PATH from .env — set
#     it literally here too; Quadlet does not interpolate .env into Volume=). ---
Volume=%h/containers/qbittorrent/config:/config
Volume=/data/downloads:/downloads

[Service]
Restart=always

[Install]
WantedBy=default.target
