# containers/qbittorrent/qbittorrent.container
# qBittorrent — joins the shared gluetun.pod (containers/gluetun/), riding its
# VPN-routed network + user namespace. No network stack of its own: if gluetun
# stops or the tunnel drops, qBittorrent loses its route (the kill switch).
# Its WebUI port is published on the POD (gluetun.pod), not here.
[Unit]
Description=qBittorrent (in gluetun.pod, riding the VPN)
# Start after gluetun and require it — the pod's route depends on it existing.
Requires=gluetun.service
After=gluetun.service
# Generous restart budget: if gluetun's container is torn down (a real restart,
# an image update, a VPN reconnect), Podman removes qBittorrent's container
# along with it (its netns just vanished) — that's the kill switch working as
# intended, not a crash. Restart=always (below) brings it back once gluetun is
# up again, but the systemd default (5 restarts / 10s) can be exhausted by a
# couple of quick back-to-back gluetun restarts, leaving qBittorrent sitting
# `inactive (dead)` until a manual `systemctl --user start qbittorrent`. Widen
# the window so it keeps retrying through that.
StartLimitIntervalSec=120
StartLimitBurst=10

[Container]
Image=lscr.io/linuxserver/qbittorrent:latest
ContainerName=qBittorrent
AutoUpdate=registry
Pod=gluetun.pod

# --- security: linuxserver/s6 image; caps intentionally NOT dropped (unverified
#     against this image's init, like the Docker stack). NoNewPrivileges kept. ---
NoNewPrivileges=true
PodmanArgs=--memory=1g --pids-limit=400

# --- UID mapping: deliberately NO keep-id anywhere in this pod (verified —
#     see gluetun.pod's comment). gluetun.pod has no UserNS= line, so this
#     container uses Podman's default rootless mapping: your real host UID
#     1000 is NOT the same identity as container UID 1000 (PUID=1000 below) —
#     it maps to a subuid-derived host UID instead (e.g. 100999). New files
#     qBittorrent creates are self-consistent (it can always read/write its
#     own), but PRE-EXISTING files owned by your real UID 1000 (anything you
#     had before this container touched it — migrated downloads, a fresh
#     bind-mounted config dir) need a ONE-TIME ownership fix so the container
#     can see them as "its own":
#       podman unshare chown -R 1000:1000 ~/containers/qbittorrent/config
#       podman unshare chown -R 1000:1000 /data/downloads
#     This re-points real on-disk ownership to whatever host UID the
#     container resolves "1000:1000" to — after running it, those files won't
#     show as owned by "hussein" in a normal `ls -la` anymore; use
#     `podman unshare <cmd>` to inspect/manage them directly from the host
#     when needed. This trade-off was chosen deliberately over UserNS=keep-id:
#     gluetun's nftables kill-switch has no documented way to disable it and
#     hard-fails under keep-id's fragmented UID mapping (confirmed live), so
#     keep-id was dropped pod-wide rather than break the VPN tunnel.
EnvironmentFile=%h/containers/qbittorrent/.env
Environment=WEBUI_PORT=8080

# --- storage: app config + your downloads disk (DOWNLOADS_PATH — set it
#     literally here too; Quadlet does not interpolate .env into Volume=). ---
Volume=%h/containers/qbittorrent/config:/config
Volume=/data/downloads:/downloads

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
