# containers/atvloadly/scripts/atvloadly-refresh.timer
[Unit]
Description=Trigger atvloadly app refresh daily at 20:30

[Timer]
OnCalendar=*-*-* 20:30:00
Persistent=true

[Install]
WantedBy=timers.target
